# sentry.service
[Unit]
Description=Sentry service

[Service]
ExecStart=/usr/bin/supervisord -c /etc/sentry/supervisord.conf
User=sentry
Group=sentry
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
